# deploy/lite/supervisord.conf
; =============================================================================
; Vexa Lite (v0.12) — single-container process supervisor
; =============================================================================
; Runs the v0.12 control plane in ONE container: the four long-running Python services
; (admin-api · runtime · meeting-api · gateway) + the agent control plane (agent-api) +
; over an internal redis + a shared Xvfb/PulseAudio audio-video stack.
;
; The bot and the agent WORKER are NOT supervised here: the runtime (RUNTIME_BACKEND=process)
; spawns them as CHILD PROCESSES on demand via the in-container launchers
; (BOT_COMMAND / AGENT_WORKER_COMMAND — defaulted by the entrypoint to the /usr/local/bin
; launchers, interpolated here via %(ENV_…)s so operator-provided container env wins),
; sharing this container's display/audio — the lite analog of compose's docker-socket spawn.
; =============================================================================
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info

; supervisorctl wiring — without these sections `supervisorctl` (the README's debug entrypoint) has
; no socket to talk to. This config IS the main file (run directly), so it carries them itself.
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ─── Internal Redis (bus + scheduler + per-dispatch streams) ─────────────────────────────────────
[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --port 6379 --dir /var/lib/redis --appendonly yes --appendfsync everysec --protected-mode yes
autostart=true
autorestart=true
priority=5
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; ─── Xvfb — virtual X11 display for the headful Chromium the bot drives ───────────────────────────
[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=10
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; ─── fluxbox — a WM so the browser window manages cleanly ─────────────────────────────────────────
[program:fluxbox]
command=fluxbox
environment=DISPLAY=":99"
autostart=true
autorestart=true
priority=11
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; ─── PulseAudio — the audio graph for the bot's capture/speak path ────────────────────────────────
[program:pulseaudio]
command=pulseaudio --system --disallow-exit --disallow-module-loading=false --log-level=notice
autostart=true
autorestart=true
priority=12
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; ─── PulseAudio sink setup (once, after PulseAudio is ready) ──────────────────────────────────────
[program:pulseaudio-setup]
command=/usr/local/bin/setup-pulseaudio-sinks.sh
autostart=true
autorestart=false
startsecs=0
priority=13
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; ─── x11vnc — share the Xvfb display (browser view, debugging) ────────────────────────────────────
[program:x11vnc]
command=x11vnc -display :99 -forever -nopw -shared -localhost -rfbport 5900 -q
autostart=true
autorestart=true
priority=14
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; ─── websockify — bridge VNC (5900) → WebSocket (6080) for noVNC ──────────────────────────────────
[program:websockify]
command=websockify --web /usr/share/novnc 6080 localhost:5900
autostart=true
autorestart=true
priority=15
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; ─── admin-api (:8001) — identity: users/tokens + /internal/validate (the gateway's authz oracle) ──
[program:admin-api]
command=/opt/venvs/admin/bin/python -m admin_api
directory=/app/admin-api
autostart=true
autorestart=true
priority=30
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/admin-api/src",HOST="127.0.0.1",PORT="8001",LOG_LEVEL="%(ENV_LOG_LEVEL)s",DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s",ADMIN_API_TOKEN="%(ENV_ADMIN_API_TOKEN)s",INTERNAL_API_SECRET="%(ENV_INTERNAL_API_SECRET)s"

; ─── runtime (:8090) — runtime.v1; spawns bot + agent-worker as CHILD PROCESSES (process backend) ──
[program:runtime]
command=/opt/venvs/runtime/bin/python -m runtime_kernel
directory=/app/runtime
autostart=true
autorestart=true
priority=35
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/runtime/src",HOST="127.0.0.1",PORT="8090",LOG_LEVEL="%(ENV_LOG_LEVEL)s",RUNTIME_BACKEND="process",REDIS_URL="%(ENV_REDIS_URL)s",DISPLAY=":99",BOT_COMMAND="%(ENV_BOT_COMMAND)s",AGENT_WORKER_COMMAND="%(ENV_AGENT_WORKER_COMMAND)s",BROWSER_IMAGE="lite-process",AGENT_IMAGE="lite-process",INTERNAL_API_SECRET="%(ENV_INTERNAL_API_SECRET)s",HOST_CLAUDE_CREDENTIALS="%(ENV_HOST_CLAUDE_CREDENTIALS)s",ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",ANTHROPIC_AUTH_TOKEN="%(ENV_ANTHROPIC_AUTH_TOKEN)s",ANTHROPIC_BASE_URL="%(ENV_ANTHROPIC_BASE_URL)s",ANTHROPIC_MODEL="%(ENV_ANTHROPIC_MODEL)s",ANTHROPIC_DEFAULT_OPUS_MODEL="%(ENV_ANTHROPIC_DEFAULT_OPUS_MODEL)s",ANTHROPIC_DEFAULT_SONNET_MODEL="%(ENV_ANTHROPIC_DEFAULT_SONNET_MODEL)s",ANTHROPIC_DEFAULT_HAIKU_MODEL="%(ENV_ANTHROPIC_DEFAULT_HAIKU_MODEL)s"

; ─── agent-api (:8100) — the unit control plane: /invocations · /api/chat (SSE) · /api/sessions ────
; LITE OVERRIDES: every VEXA_*_API_URL default in agent config points at a COMPOSE service hostname
; (runtime-api/agent-api/meeting-api), which does NOT resolve in the single-container lite deploy —
; each must be overridden to localhost here. Missing VEXA_MEETING_API_URL made the live-transcript
; owner check (GET meeting-api/meetings/{id}) throw on DNS → fail-closed → EVERY live SSE 403'd, so
; the terminal panel stayed blank on lite though transcription itself worked (#585).
[program:agent-api]
command=/opt/venvs/agent/bin/uvicorn control_plane.api:app --host 0.0.0.0 --port 8100 --log-level %(ENV_LOG_LEVEL)s
directory=/app/agent
autostart=true
autorestart=true
priority=45
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/agent",VEXA_AGENT_API_PORT="8100",VEXA_RUNTIME_API_URL="http://localhost:8090",VEXA_AGENT_API_SELF_URL="http://localhost:8100",VEXA_MEETING_API_URL="http://localhost:8080",VEXA_AGENT_PROFILE="agent",VEXA_WORKSPACES_DIR="/workspaces",VEXA_WORKSPACE_SEEDS_DIR="/app/agent/workspace-seeds",VEXA_TOOLS_SEED_DIR="/app/agent/tools-seed",VEXA_AGENT_DEFAULT_SUBJECT="%(ENV_VEXA_AGENT_DEFAULT_SUBJECT)s",VEXA_REDIS_URL="%(ENV_REDIS_URL)s",VEXA_DISPATCH_SIGNING_KEY="%(ENV_VEXA_DISPATCH_SIGNING_KEY)s",VEXA_GATEWAY_URL="http://localhost:8056",VEXA_BOT_API_KEY="%(ENV_VEXA_BOT_API_KEY)s",VEXA_LOG_LEVEL="%(ENV_LOG_LEVEL)s",VEXA_AGENT_MODEL="%(ENV_VEXA_AGENT_MODEL)s",VEXA_MEETING_MODEL="%(ENV_VEXA_MEETING_MODEL)s",ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",ANTHROPIC_AUTH_TOKEN="%(ENV_ANTHROPIC_AUTH_TOKEN)s",ANTHROPIC_BASE_URL="%(ENV_ANTHROPIC_BASE_URL)s",ANTHROPIC_MODEL="%(ENV_ANTHROPIC_MODEL)s",ANTHROPIC_DEFAULT_OPUS_MODEL="%(ENV_ANTHROPIC_DEFAULT_OPUS_MODEL)s",ANTHROPIC_DEFAULT_SONNET_MODEL="%(ENV_ANTHROPIC_DEFAULT_SONNET_MODEL)s",ANTHROPIC_DEFAULT_HAIKU_MODEL="%(ENV_ANTHROPIC_DEFAULT_HAIKU_MODEL)s"

; ─── meeting-api (:8080) — lifecycle · bot_spawn · collector · recordings (→ minio) ───────────────
[program:meeting-api]
command=/opt/venvs/meeting/bin/python -m meeting_api
directory=/app/meeting-api
autostart=true
autorestart=true
priority=40
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/meeting-api/src",HOST="127.0.0.1",PORT="8080",LOG_LEVEL="%(ENV_LOG_LEVEL)s",DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s",REDIS_URL="%(ENV_REDIS_URL)s",RUNTIME_API_URL="http://localhost:8090",MEETING_API_URL="http://localhost:8080",INTERNAL_API_SECRET="%(ENV_INTERNAL_API_SECRET)s",ADMIN_TOKEN="%(ENV_ADMIN_API_TOKEN)s",TRANSCRIPTION_SERVICE_URL="%(ENV_TRANSCRIPTION_SERVICE_URL)s",TRANSCRIPTION_SERVICE_TOKEN="%(ENV_TRANSCRIPTION_SERVICE_TOKEN)s",MINIO_ENDPOINT="%(ENV_MINIO_ENDPOINT)s",MINIO_ACCESS_KEY="%(ENV_MINIO_ACCESS_KEY)s",MINIO_SECRET_KEY="%(ENV_MINIO_SECRET_KEY)s",MINIO_BUCKET="%(ENV_MINIO_BUCKET)s",MINIO_SECURE="%(ENV_MINIO_SECURE)s"

; ─── gateway (:8056) — the one front door: auth · routing · /ws fan-out ───────────────────────────
[program:gateway]
command=/opt/venvs/gateway/bin/python -m gateway
directory=/app/gateway
autostart=true
autorestart=true
priority=50
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
; fastapi-guard edge layer (edge_guard.py): ON by default with generous limits (owner ruling);
; opt out with -e GUARD_ENABLED=false (entrypoint.sh defaults it, then it passes through here).
; GUARD_WS_ENABLED is the optional /ws hook (off, opt-in): in-memory per-process — does NOT share
; ban/rate-limit state with the Redis-backed HTTP layer, and under uvicorn --workers N>1 each
; worker keeps an independent WS ban set. With GUARD_ENABLE_REDIS=false AND uvicorn --workers N>1
; the HTTP rate-limit buckets and auto-bans are likewise per-process: the effective limit becomes
; N × GUARD_RATE_LIMIT_RPM and bans do not propagate across workers. Other GUARD_* tuning keys inherit from container env.
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/gateway/src",HOST="0.0.0.0",PORT="8056",LOG_LEVEL="%(ENV_LOG_LEVEL)s",ADMIN_API_URL="http://localhost:8001",MEETING_API_URL="http://localhost:8080",AGENT_API_URL="http://localhost:8100",REDIS_URL="%(ENV_REDIS_URL)s",INTERNAL_API_SECRET="%(ENV_INTERNAL_API_SECRET)s",GUARD_ENABLED="%(ENV_GUARD_ENABLED)s",GUARD_WS_ENABLED="%(ENV_GUARD_WS_ENABLED)s"


; ─── terminal (:3001) — the agent-domain browser-CLI workbench (Next.js + custom server.mjs) ──────
[program:terminal]
; Source the self-host key file (provision-key.sh) so VEXA_API_KEY is present for zero-login; the
; proxy + /ws relay forward it as X-API-Key. Absent file ⇒ login mode.
command=sh -c 'set -a; [ -f /run/vexa/key.env ] && . /run/vexa/key.env; set +a; exec node server.mjs'
directory=/app/terminal
autostart=true
autorestart=true
priority=61
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="3001",HOST="0.0.0.0",AGENT_API_URL="http://localhost:8100",GATEWAY_URL="http://localhost:8056",VEXA_API_URL="http://localhost:8056",VEXA_ADMIN_API_URL="http://localhost:8001",VEXA_ADMIN_API_KEY="%(ENV_ADMIN_API_TOKEN)s",VEXA_API_KEY="%(ENV_VEXA_API_KEY)s",VEXA_BOT_API_KEY="%(ENV_VEXA_BOT_API_KEY)s",NEXTAUTH_URL="%(ENV_TERMINAL_PUBLIC_URL)s",NEXTAUTH_SECRET="%(ENV_NEXTAUTH_SECRET)s"

; ─── process group ────────────────────────────────────────────────────────────────────────────────
[group:vexa]
programs=redis,xvfb,fluxbox,pulseaudio,pulseaudio-setup,x11vnc,websockify,admin-api,runtime,agent-api,meeting-api,gateway,terminal
priority=999
